# supervisord/datacoin.conf
[program:datacoin]
environment=HOME="/home/datacoin"
command=/home/datacoin/datacoin-hp/src/datacoind
directory=/home/datacoin
user=datacoin

autostart=true
autorestart=true
redirect_stderr=true
