# supervisord/metadisk-cloudsync.conf
[program:metadisk-cloudsync]
command=/home/metadisk/web-core/.env/bin/python cloudsync.py
directory=/home/metadisk/web-core
user=metadisk

autostart=true
autorestart=true
redirect_stderr=true
