# chef/site-cookbooks/storj/templates/default/supervisord-metadisk-cloudsync.conf.erb
[program:metadisk-cloudsync]
command=/home/<%= @user %>/web-core/.env/bin/python cloudsync.py
directory=/home/<%= @user %>/web-core
user=<%= @user %>

autostart=true
autorestart=true
redirect_stderr=true
