# chef/site-cookbooks/storj/templates/default/supervisord-websockets.conf.erb
[program:websockets]
command=/home/<%= @user %>/metadisk-websockets/.env/bin/python index.py
directory=/home/<%= @user %>/metadisk-websockets
user=<%= @user %>

autostart=true
autorestart=true
redirect_stderr=true